# monitoring-stack/systemd/grafana-pf.service
# If your kubectl path differs, change /usr/local/bin/kubectl.

[Unit]
Description=Grafana port-forward (K3s)
After=network-online.target k3s.service
Wants=network-online.target

[Service]
Environment=KUBECONFIG=/etc/rancher/k3s/k3s.yaml
Environment=PATH=/usr/local/bin:/usr/bin:/bin
SyslogIdentifier=grafana-pf
ExecStart=/usr/local/bin/kubectl -n monitoring port-forward svc/prometheus-stack-grafana 3000:80 --address 127.0.0.1
Restart=always
RestartSec=2

[Install]
WantedBy=multi-user.target